# services/mysql-backup/mysql-backup.timer
[Unit]
Description=Run MYSQL Backup every 1 hour

[Timer]
OnCalendar=hourly
Unit=mysql-backup.service

[Install]
WantedBy=timers.target